# scripts/com.kemmer.df-heylou-openai-extension.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- df-heylou-openai-extension macOS LaunchAgent [CRUX-MK]
     Welle-39 Foundation-DF: HeyLou-Sub-Funktion-Extension Gemini Function-Calling. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kemmer.df-heylou-openai-extension</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/Users/make/Projects/dark-factories/df-heylou-openai-extension/scripts/run-df-heylou-openai-extension-mac.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>7200</integer>

    <key>StandardOutPath</key>
    <string>/Users/make/Projects/dark-factories/df-heylou-openai-extension/runs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/make/Projects/dark-factories/df-heylou-openai-extension/runs/stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/make/Projects/dark-factories/df-heylou-openai-extension</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>DF_HEYLOU_OPENAI_EXT_ENABLED</key>
        <string>false</string>
        <key>DF_HEYLOU_OPENAI_TENANT_ID</key>
        <string>hildesheim</string>
        <key>DF_HEYLOU_OPENAI_TIMEOUT_S</key>
        <string>30</string>
    </dict>
</dict>
</plist>
